# com.hangewubi.ci-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hangewubi.ci-watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/linhan/startup/函戈五笔/ci-watch.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/linhan/startup/函戈五笔</string>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/linhan/startup/函戈五笔/ci-logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/linhan/startup/函戈五笔/ci-logs/launchd.log</string>
</dict>
</plist>
